# unicorn-transcoder/rootfs/etc/services.d/unicorn-transcoder/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid plex npm start --prefix=/opt/UnicornTranscoder
