# unicorn-transcoder/rootfs/etc/services.d/unicorn-transcoder/run
#!/usr/bin/with-contenv sh

function runTranscoder(){
  npm start --prefix=/opt/UnicornTranscoder
}

exec s6-setuidgid plex runTranscoder
